# initfiles/bash/etc/systemd/system/hpcc-init.service.in
[Unit]
Description=hpcc-systems platform

[Service]
Type=oneshot
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=no
ExecStart=@ADMIN_PATH@/generate-hpccsystems-target
ExecStop=@ADMIN_PATH@/generate-hpccsystems-target
WorkingDirectory=@RUNTIME_PATH@
StandardOutput=journal+console
